# scripts/launchd/com.mitchell.career-ops.signal-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.mitchell.career-ops.signal-monitor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/mitchellwilliams/.nvm/versions/node/v24.14.0/bin/node</string>
    <string>/Users/mitchellwilliams/Documents/career-ops/signal-monitor.mjs</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/mitchellwilliams/Documents/career-ops</string>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/signal-monitor.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/mitchellwilliams/Documents/career-ops/data/logs/signal-monitor.log</string>
</dict>
</plist>
